# roles/base/templates/rc.local/rc.local.j2
#!/bin/sh -e
#
# {{ ansible_managed }} [ansible-playbooks v{{ ansibleplaybooks_version }}]
#
# rc.local

hostname "{{ ansible_hostname }}"
domainname "{{ ansible_domain }}"

exit 0
